# com.opentidy.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  OpenTidy LaunchAgent template.
  install.sh copies this to ~/Library/LaunchAgents/ and replaces __INSTALL_DIR__, __NODE__, __HOME__.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.opentidy.agent</string>

	<key>ProgramArguments</key>
	<array>
		<string>__NODE__</string>
		<string>__INSTALL_DIR__/apps/backend/dist/cli.js</string>
		<string>start</string>
	</array>

	<key>WorkingDirectory</key>
	<string>__INSTALL_DIR__</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>__HOME__/.local/bin:/opt/homebrew/opt/node@22/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>NODE_ENV</key>
		<string>production</string>
		<key>HOME</key>
		<string>__HOME__</string>
	</dict>

	<key>KeepAlive</key>
	<true/>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/opentidy-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/opentidy-stderr.log</string>
</dict>
</plist>
